# raspi-sump/cron/crontab
45 * * * * * /home/pi/raspi-sump/move_file.sh
