# conf/systemd/raspisump.service
[Unit]
Description=RaspiSump Service

[Service]
ExecStart=/opt/raspi-sump/bin/rsump.py
Restart=always

[Install]
WantedBy=default.target
