# conf/systemd/raspisump.service
[Unit]
Description=RaspiSump Service

[Service]
ExecStart=/usr/local/bin/rsump.py
Restart=always

[Install]
WantedBy=multi-user.target
