# conf/systemd/rsumpwebchart.service
[Unit]
Description=RaspiSump Webchart Service

[Service]
ExecStart=/opt/raspi-sump/bin/rsumpwebchart.py

[Install]
WantedBy=multi-user.target
