# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create raspisump system user if not present
        if ! getent passwd raspisump > /dev/null 2>&1; then
            adduser --system --group --no-create-home \
                    --home /var/lib/raspi-sump \
                    --gecos "Raspi-Sump system user" \
                    raspisump
        fi

        # Add raspisump to gpio group
        if getent group gpio > /dev/null 2>&1; then
            usermod -aG gpio raspisump
        fi

        # Add the user who invoked sudo/dpkg to the raspisump group so they
        # can run manual readings without switching to the service account.
        INSTALL_USER="${SUDO_USER:-}"
        if [ -n "$INSTALL_USER" ] && getent passwd "$INSTALL_USER" > /dev/null 2>&1; then
            usermod -aG raspisump "$INSTALL_USER"
        fi

        # Create state subdirectories (systemd StateDirectory only creates the top level)
        install -d -o raspisump -g raspisump -m 0770 /var/lib/raspi-sump
        install -d -o raspisump -g raspisump -m 0770 /var/lib/raspi-sump/csv
        install -d -o raspisump -g raspisump -m 0775 /var/lib/raspi-sump/charts
        install -d -o raspisump -g raspisump -m 0770 /var/lib/raspi-sump/support

        # Create log directory (systemd LogsDirectory may not be present before first start)
        install -d -o raspisump -g raspisump -m 0770 /var/log/raspi-sump

        # Fix ownership and permissions recursively — catches any subdirs (e.g. charts/year/month)
        # created by an older install with the wrong group or mode.
        chown -R raspisump:raspisump \
            /var/lib/raspi-sump/csv \
            /var/lib/raspi-sump/charts \
            /var/lib/raspi-sump/support \
            /var/log/raspi-sump
        find /var/lib/raspi-sump/csv \
             /var/lib/raspi-sump/support \
             /var/log/raspi-sump \
             -type d -exec chmod g+rwxs {} \;
        find /var/lib/raspi-sump/charts \
             -type d -exec chmod g+rwxs,o+rx {} \;

        # Set setgid bit on top-level shared dirs so new files/dirs inherit group 'raspisump'.
        chmod g+s \
            /var/lib/raspi-sump/csv \
            /var/lib/raspi-sump/support \
            /var/log/raspi-sump
        # Charts dir is world-readable so nginx can serve files without group membership.
        chmod g+s,o+rx /var/lib/raspi-sump/charts

        # Install config files from examples only on first install (never overwrite on upgrade)
        install -d -o root -g root -m 0755 /etc/raspi-sump
        if [ ! -f /etc/raspi-sump/raspisump.conf ]; then
            cp /usr/share/raspi-sump/examples/raspisump.conf /etc/raspi-sump/raspisump.conf
            chmod 0644 /etc/raspi-sump/raspisump.conf
        fi
        if [ ! -f /etc/raspi-sump/credentials.conf ]; then
            cp /usr/share/raspi-sump/examples/credentials.conf /etc/raspi-sump/credentials.conf
        fi

        # Lock down credentials file — only root and raspisump user can read it
        chown root:raspisump /etc/raspi-sump/credentials.conf
        chmod 0640 /etc/raspi-sump/credentials.conf

        # Write version string into the web footer file
        PKG_VERSION=$(dpkg-query --showformat='${Version}' --show raspisump 2>/dev/null || echo "unknown")
        echo "Raspisump Version ${PKG_VERSION} | Copyright 2026 Al Audet" \
            > /usr/share/raspi-sump/web/css/inc/VERSION

        # nginx site handling
        NGINX_SITE_ENABLED="/etc/nginx/sites-enabled/raspi-sump"

        if [ -z "$2" ]; then
            # Fresh install — check for port 80 conflicts before enabling
            NGINX_CONFLICT_MSG=""
            if systemctl is-active --quiet apache2 2>/dev/null; then
                NGINX_CONFLICT_MSG="apache2 is active"
            elif systemctl is-active --quiet lighttpd 2>/dev/null; then
                NGINX_CONFLICT_MSG="lighttpd is active"
            elif [ -L /etc/nginx/sites-enabled/default ]; then
                NGINX_CONFLICT_MSG="the nginx default site is enabled (conflicts on port 80)"
            fi

            if [ -z "$NGINX_CONFLICT_MSG" ]; then
                ln -sf /etc/nginx/sites-available/raspi-sump "$NGINX_SITE_ENABLED"
                systemctl reload nginx > /dev/null 2>&1 || true
                NGINX_ENABLED=1
            else
                NGINX_ENABLED=0
            fi
        else
            # Upgrade — respect whatever state the site is currently in
            if [ -L "$NGINX_SITE_ENABLED" ]; then
                systemctl reload nginx > /dev/null 2>&1 || true
            fi
        fi

        systemctl daemon-reload || true
        systemctl enable raspisump.service || true
        systemctl enable rsumpwebchart.timer || true

        # On upgrade, restart any services that are enabled (prerm stops them before postinst runs)
        if [ -n "$2" ]; then
            if systemctl is-enabled --quiet raspisump.service; then
                systemctl restart raspisump.service || true
            fi
            if systemctl is-enabled --quiet rsumpwebchart.timer; then
                systemctl restart rsumpwebchart.timer || true
            fi
        fi

        if [ -z "$2" ]; then
            # Fresh install
            echo ""
            echo "raspisump has been installed and the service has been enabled."
            echo ""
            echo "Before starting the service, edit the configuration files:"
            echo "  /etc/raspi-sump/raspisump.conf"
            echo "  /etc/raspi-sump/credentials.conf"
            echo ""
            echo "Then start the service with:"
            echo "  sudo systemctl start raspisump.service"
            echo ""
            if [ "$NGINX_ENABLED" -eq 1 ]; then
                echo "The raspisump web interface is available at http://<your-pi-ip>/"
                echo ""
            else
                echo "Web server conflict detected: ${NGINX_CONFLICT_MSG}."
                echo "The raspisump nginx config has been installed but not activated:"
                echo "  /etc/nginx/sites-available/raspi-sump"
                echo "Resolve the conflict and enable it manually:"
                echo "  sudo ln -s /etc/nginx/sites-available/raspi-sump /etc/nginx/sites-enabled/raspi-sump"
                echo "  sudo systemctl reload nginx"
                echo ""
            fi
        else
            # Upgrade
            echo ""
            echo "raspisump has been upgraded from $2."
            echo ""
            echo "Your configuration files have NOT been modified:"
            echo "  /etc/raspi-sump/raspisump.conf"
            echo "  /etc/raspi-sump/credentials.conf"
            echo ""
            echo "Updated example configs are available at:"
            echo "  /usr/share/raspi-sump/examples/raspisump.conf"
            echo "  /usr/share/raspi-sump/examples/credentials.conf"
            echo ""
            echo "Review the examples for any new settings before restarting the service."
            echo ""
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#
exit 0
